# deploy/deskbrid.service
[Unit]
Description=Deskbrid desktop agent HAL
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/deskbrid daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
